# rest/src/deb_template/init/docussandra.conf
description "Docussandra server version VERSION"

start on runlevel [2345]
stop on runlevel [!2345]

env USER=appuser

script
    echo $$ > /var/run/docussandraVERSION.pid
    . /etc/environment
    export WW_DNS
    touch /var/log/docussandra/docussandra.start
    chmod 666 /var/log/docussandra/docussandra.start
    su - $USER -c "echo `date` \"Starting up Docussandra version: VERSION and with properties: $WW_DNS\""
    su - $USER -c "echo `date` \"Starting up Docussandra version: VERSION and with properties: $WW_DNS\" >> /var/log/docussandra/docussandra.start"
    su - $USER -c "cd /home/appuser/docussandra-rest-VERSION; java -jar docussandra-rest-jar-with-dependencies.jar $WW_DNS &> /dev/null 2>&1"
end script

post-stop script
    rm -f /var/run/docussandraVERSION.pid
    sleep 2
end script
